# asterisk-cc-phase1/ops/systemd/cc-gsm-routes.service
[Unit]
Description=Static routes to GSM (10.1.5.0/24 + /29 + /27) via 172.16.4.1
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=GSM_ROUTE_NET=10.1.5.0/24
Environment=GSM_ROUTE_SIGNAL_NET=10.1.5.8/29
Environment=GSM_ROUTE_MEDIA_NET=10.1.5.64/27
Environment=GSM_ROUTE_VIA=172.16.4.1
Environment=GSM_ROUTE_DEV=enp13s4f0
ExecStart=/bin/bash /opt/call-center/asterisk-cc-phase1/scripts/apply-gsm-routes.sh
ExecStop=/bin/sh -c 'ip route del 10.1.5.0/24 via 172.16.4.1 dev enp13s4f0 2>/dev/null; ip route del 10.1.5.8/29 via 172.16.4.1 dev enp13s4f0 2>/dev/null; ip route del 10.1.5.64/27 via 172.16.4.1 dev enp13s4f0 2>/dev/null; true'

[Install]
WantedBy=multi-user.target
