# asterisk-cc-phase1/ops/systemd/native/cc-media-debug.service
[Unit]
Description=Asterisk CC media debug settings
After=cc-asterisk.service
Requires=cc-asterisk.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/cc/cc.env
ExecStart=/bin/bash /opt/cc/scripts/cc-media-debug-apply.sh

[Install]
WantedBy=multi-user.target
